# newsfeed.service
[Unit]
Description=Newsfeed — hourly news digest check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/newsfeed/venv/bin/python3 /opt/newsfeed/runner.py
